# scripts/pi-bot.service
# systemd unit for pi-Bot.
#
# Install:
#   sudo cp scripts/pi-bot.service /etc/systemd/system/pi-bot.service
#   sudo systemctl edit pi-bot.service           # adjust User/WorkingDirectory
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pi-bot
#
# Logs:
#   journalctl -u pi-bot -f
#   tail -f <WorkingDirectory>/logs/app.log
#
# systemd sends SIGTERM on `stop`/`restart`, which the bot handles cleanly via
# its asyncio shutdown path (final state save before exit).

[Unit]
Description=pi-Bot Discord dice rolling bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust these to your deployment.
User=pibot
Group=pibot
WorkingDirectory=/opt/pi-bot
EnvironmentFile=/opt/pi-bot/.env
ExecStart=/opt/pi-bot/.venv/bin/python -m pi_bot

# Sane restart policy. Don't restart on a clean exit (e.g. owner-issued
# shutdown) but do restart on crashes and signals.
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60
StartLimitBurst=5

# Allow up to 30s for the bot to flush state on shutdown.
TimeoutStopSec=30

# Hardening — adjust if you deploy elsewhere.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pi-bot/data /opt/pi-bot/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
