# scripts/sirrmizan.service
# systemd unit for SirrMizan.
#
# Install:
#   sudo cp scripts/sirrmizan.service /etc/systemd/system/sirrmizan.service
#   sudo systemctl edit sirrmizan.service           # adjust User/WorkingDirectory
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now sirrmizan
#
# Logs:
#   journalctl -u sirrmizan -f
#   tail -f <WorkingDirectory>/logs/app.log
#
# systemd sends SIGTERM on `stop`/`restart`, which the bot handles cleanly via
# its asyncio shutdown path (final state save before exit).

[Unit]
Description=SirrMizan Discord dice rolling bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust these to your deployment.
User=sirrmizan
Group=sirrmizan
WorkingDirectory=/opt/sirrmizan
EnvironmentFile=/opt/sirrmizan/.env
ExecStart=/opt/sirrmizan/.venv/bin/python -m sirrmizan

# Sane restart policy. Don't restart on a clean exit (e.g. owner-issued
# shutdown) but do restart on crashes and signals.
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60
StartLimitBurst=5

# Allow up to 30s for the bot to flush state on shutdown.
TimeoutStopSec=30

# Hardening — adjust if you deploy elsewhere.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sirrmizan/data /opt/sirrmizan/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
